# roles/poll/files/poll.service
[Unit]
Description=Poll web client

[Service]
User=root
WorkingDirectory=/poll
ExecStart=/usr/bin/python3 wsgi.py
Restart=always
